# host/audit_host.cron
* * * * * root PATH=/sbin:/usr/sbin:/bin:/usr/bin:${PATH};/misc/scratch/audit_host.bin
